# etc/supervisord.d/grafana.ini
[program:grafana3]
directory=/usr/share/grafana/
command=/opt/qnib/grafana3/bin/start.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
